# deploy/systemd/power-detector.service
[Unit]
Description=Power Detector Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/power-detector
EnvironmentFile=-/opt/power-detector/.env
ExecStart=/opt/power-detector/.venv/bin/python /opt/power-detector/power_detector.py --config /opt/power-detector/config.yaml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
